# integration-tests/docker/service-supervisords/historical.conf
[program:druid-historical]
command=java
  -server
  -Xmx512m
  -Xms512m
  -XX:NewSize=256m
  -XX:MaxNewSize=256m
  -XX:+UseG1GC
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Dlog4j.configurationFile=/shared/docker/lib/log4j2.xml
  -cp %(ENV_DRUID_COMMON_CONF_DIR)s:%(ENV_DRUID_SERVICE_CONF_DIR)s:/shared/docker/lib/*:/usr/local/druid/lib/mysql-connector-java.jar
  org.apache.druid.cli.Main server historical
redirect_stderr=true
priority=100
autorestart=false
stdout_logfile=/shared/logs/historical.log
environment=AWS_REGION=us-east-1
