# integration-tests/docker/service-supervisords/kafka.conf
[program:kafka]
command=/usr/local/kafka/bin/kafka-server-start.sh /usr/local/kafka/config/server.properties
priority=0
stdout_logfile=/shared/logs/kafka.log
