# integration-tests/docker/service-supervisords/metadata-storage.conf
[program:mysql]
command=/run-mysql.sh
user=mysql
priority=0
stdout_logfile=/shared/logs/mysql.log
