# integration-tests/docker/service-supervisords/zookeeper.conf
[program:zookeeper]
command=/usr/local/zookeeper-%(ENV_ZK_VERSION)s/bin/zkServer.sh start-foreground
user=daemon
priority=0
stdout_logfile=/shared/logs/zookeeper.log
