# integration-tests/docker/mysql.conf
[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
         --bind-address=0.0.0.0
user=mysql
priority=0
stdout_logfile=/shared/logs/mysql.log
